# root/etc/services.d/inotify_modify/run
#!/usr/bin/with-contenv bash

# starting inotify to watch /config/ddclient.conf and restart ddclient if changed.
while inotifywait -e close_write /config/ddclient.conf; do
	cp /config/ddclient.conf /ddclient.conf
	chmod 600 /ddclient.conf
	chown ddclient:ddclient /ddclient.conf
	kill "$(cat /var/run/ddclient/ddclient.pid)"
	echo "ddclient have been restarted"
done