# root/etc/services.d/ddclient/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc /usr/bin/ddclient -foreground -file /ddclient.conf
